# misc/ovenmediaengine.service
[Unit]
Description=OvenMediaEngine
After=network-online.target

[Service]
Type=forking
PIDFile=/var/run/ovenmediaengine.pid
WorkingDirectory=/usr/share/ovenmediaengine
ExecStart=/usr/share/ovenmediaengine/ome_launcher.sh -d
Restart=always
RestartSec=2
RestartPreventExitStatus=1
StandardOutput=null
StandardError=null
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target